# systemd-units/change-wallpaper.service
[Unit]
Description=Change the wallpaper every day.

[Service]
Environment="PATH=/usr/bin:/bin:%h/.local/bin"
Type=simple
ExecStart=/bin/sh %h/.fehbg
